# systemd/tavazon.service
[Unit]
Description=Tavazon - upload/download ratio balancer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tavazon
ExecStart=/opt/tavazon/tavazon -config /opt/tavazon/config.json
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/tavazon
ProtectHome=true

[Install]
WantedBy=multi-user.target
